# utils/cronFiles/crontab
* * * * * root /bin/date +\%d/\%m/\%Y-\%H:\%M:\%S >> /home/crontestOutput
* * * * * root /bin/sh /backup/backup.sh



