# utils/cronFiles/crontab
* * * * * root /bin/sh /home/backup.sh
* * * * * root /bin/date +\%d/\%m/\%Y-\%H:\%M:\%S >> /home/crontestOutput



